# systemd/cv-analytics.service
[Unit]
Description=CANNECT.AI CV Analytics Pipeline
Documentation=https://cannect.ai/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cannect
Group=cannect
WorkingDirectory=/opt/cannect/cv-analytics
ExecStart=/opt/cannect/cv-analytics/.venv/bin/python -m src.main
Restart=always
RestartSec=5
WatchdogSec=120

# Environment
EnvironmentFile=/opt/cannect/cv-analytics/.env

# Resource limits
MemoryMax=1G
CPUQuota=80%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/cannect/cv-analytics/buffer
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cv-analytics

[Install]
WantedBy=multi-user.target
